# provisioning/roles/apps/parrotmail/templates/parrotmail.upstart
description     "mailserver inbox app"
author          "th0ughtcr1me <th0ughtcr1me@protonmail.ch>"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 5 60

chdir /srv/mail

setuid {{ parrotmail_user }}
setgit {{ parrotmail_group }}
env INBOX_PORT={{ parrotmail_smtp_port }}
env PYTHONPATH=/srv/src
exec {{ parrotmail_venv_dir }}/bin/python {{ parrotmail_install_dir }}/mailserver.py
